# etc/systemd/system/docker.service.d/override.conf
[Service]
TasksMax=infinity
Restart=on-failure
