# etc/systemd/system/run-regularly.service
[Unit]
Description=Custom maintenance scripts

[Service]
Type=oneshot

ExecStart=/usr/local/sbin/run-regularly

[Install]
WantedBy=default.target
