# etc/systemd/system/run-regularly.timer
[Unit]
Description=Run run-regularly daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
