# etc/systemd/system/run-on-boot.service
[Unit]
Description=Custom boot scripts

[Service]
Type=oneshot

ExecStart=/home/secure/bin/run-on-boot

[Install]
WantedBy=multi-user.target

