# etc/systemd/system/run-regularly.service
[Unit]
Description=Custom maintenance scripts

[Service]
Type=oneshot

ExecStart=/home/secure/bin/run-regularly

[Install]
WantedBy=default.target
